# scripts/golden-image/firstboot-hostname.service
[Unit]
Description=Auto Assign Hostname on First Boot
Documentation=https://github.com/your-org/k8s-node-image
After=network-online.target
Wants=network-online.target
Requires=firstboot-reset.service
After=firstboot-reset.service

[Service]
Type=oneshot
ExecStart=/opt/golden-image/firstboot-hostname.sh
RemainAfterExit=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

