# scripts/golden-image/kubeadm-join.service
[Unit]
Description=Auto Join Kubernetes Cluster
Documentation=https://github.com/your-org/k8s-node-image
After=network-online.target kubelet.service containerd.service
Wants=network-online.target
Requires=firstboot-hostname.service
After=firstboot-hostname.service

[Service]
Type=oneshot
ExecStart=/opt/golden-image/kubeadm-join.sh
RemainAfterExit=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

